# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:agent]
command=python3 /app/agent.py
directory=/app
autorestart=false
stdout_logfile=/var/log/agent.log
stderr_logfile=/var/log/agent.log 